# .config/systemd/user/pipewire.service.d/override.conf
[Service]
LimitRTPRIO=95
LimitMEMLOCK=infinity
Slice=app.slice
